# systemd/homelab-backend.service
[Unit]
Description=HomeLab VM Manager - Backend (FastAPI)
After=network.target

[Service]
Type=simple
# Adaptez User et les chemins si votre installation diffère.
User=trinity
WorkingDirectory=/home/trinity/homelab-vm-manager/backend
ExecStart=/home/trinity/homelab-vm-manager/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8420
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
